# infrastructure/environments/dev/b_configuration/supervisor/photos-service.ini
[program:photos-service]
user=appuser
directory=/opt/apps/photos/
command=java -jar photos-service.jar

autostart=true
autorestart=true

stdout_logfile=syslog
stderr_logfile=syslog
