# infrastructure/environments/dev/b_configuration/supervisor/photos-jobs.ini
[program:photos-jobs]
user=appuser
directory=/opt/apps/photos/
command=java -jar photos-jobs.jar --photos-post-processor

autostart=true
autorestart=true

stdout_logfile=syslog
stderr_logfile=syslog
